# configs/airootfs/etc/profile.d/pit-aliases.sh
# Ferrari OS package manager aliases
alias update='pit -Syu'
alias install='pit -S'
alias remove='pit -R'
alias search='pit -Ss'
alias clean='pit -Sc'
alias info='pit -Si'
alias packages='pit -Q'
