# configs/airootfs/etc/profile.d/pit-aliases.sh
#!/bin/bash
# Ferrari OS package manager aliases
# System update
alias pit-update='pit -Syu'
# Package installation
alias pit-install='pit -S'
# Package removal
alias pit-remove='pit -R'
# Package search
alias pit-search='pit -Ss'
# Clean package cache
alias pit-clean='pit -Sc'
# Package information
alias pit-info='pit -Si'
# List installed packages
alias pit-list='pit -Q'
